# provisioners/roles/Snort/files/rc.local

#ethtool command to reduce processing at eth0
ethtool -K eth0 gro off lro off rx off tx off gso off

mkdir -p /opt/bro/pcaps

screen -S capture -t capture -d -m bash -c "nice -n 15 tcpdump -i eth0 -G 3600 -w '/opt/bro/pcaps/'`hostname -s`'.%Y%m%d%H%M%S.pcap' -z bzip2 '(tcp[13] & 0x17 !=0x10) or not tcp'"

/usr/bin/zeekctl deploy

/usr/local/bin/pi_show.py >/tmp/pi_show.out 2>&1 &

exit 0